# docker/services.d/postgres/run
#!/command/with-contenv sh
IMAGE_DB_PROVIDER=$(cat /etc/inker/image-db-provider 2>/dev/null || echo "unknown")
if [ "$IMAGE_DB_PROVIDER" != "postgresql" ] && [ "$DB_PROVIDER" = "postgresql" ]; then
  echo "ERROR: Container was built with DB_PROVIDER=$IMAGE_DB_PROVIDER but started with DB_PROVIDER=postgresql. PostgreSQL is not installed in this image."
  s6-svscanctl -t /run/service
  exit 1
fi

if [ "$DB_PROVIDER" != "postgresql" ]; then
  s6-svc -O .
  exit 0
fi
exec s6-setuidgid postgres /usr/lib/postgresql/17/bin/postgres -D /var/lib/postgresql/17/main
